# cron.d/md-compare
# Create a monthly rsync diff report
REPORTDIR=/var/log/reports/md-diff/
@daily root nice -n 5 rsync -aHAXni --delete /srv/current/ /srv/backup/ > $REPORTDIR/$(/usr/bin/date -I).txt
